# chatbot.service
[Unit]
Description=YC Spark Chatbot Demo with LangGraph Agent
After=network.target
Wants=network.target

[Service]
Type=simple
User={{USER}}
Group={{USER}}
WorkingDirectory={{WORKING_DIR}}
ExecStart={{WORKING_DIR}}/run_app.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PATH={{WORKING_DIR}}/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH={{WORKING_DIR}}

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target